# webhook.conf
[program:git_webhook]
command=/usr/bin/php -S 0.0.0.0:12345
directory=/root/git-webhook/public
restart=always
user=root
